# ZendServer/Homestead/php5-fpm-pid-link.conf
# Start when pc starts up
start on runlevel [2345]
# Stop when pc shuts down 
stop on runlevel [016]

# Start Script 
ln -s /usr/local/zend/tmp/php-fpm.sock /var/run/php5-fpm.sock
